# build/DEBIAN/postinst
#!/bin/sh
update-rc.d horizon defaults
systemctl daemon-reload
